# claudegram.service
[Unit]
Description=Claudegram Bot
After=network.target

[Service]
Type=simple
User=claude-bot
Group=claude-bot
WorkingDirectory=/opt/claudegram
EnvironmentFile=/opt/claudegram/.env
ExecStart=/opt/claudegram/venv/bin/python -m bot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
